# ansible/roles/s3_backup/templates/s3-backup.timer.j2
[Unit]
Description=Scheduled S3 backup of {{ s3_backup_name }}

[Timer]
OnCalendar={{ s3_backup_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
