# ansible/roles/kiosk/templates/kiosk.service.j2
{# Managed by the kiosk role. #}
[Unit]
Description=Grafana wall kiosk (cage + Chromium)
After=systemd-user-sessions.service network-online.target getty@tty1.service
Wants=network-online.target
# Take tty1 from the console getty (StandardInput=tty-fail below needs it free).
Conflicts=getty@tty1.service

[Service]
Type=simple
User={{ kiosk_user }}
# A logind (PAM) session gives cage a seat plus the DRM/input devices, so no
# separate seatd is needed. tty1 so the kiosk owns the primary console.
PAMName=login
TTYPath=/dev/tty1
StandardInput=tty-fail
StandardOutput=journal
StandardError=journal
ExecStart=/usr/bin/cage -- /usr/local/bin/kiosk-session
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
